# rc.local
#!/bin/sh -e
#LED
config-pin P8_38 1

#SWITCH 1
config-pin P8_39 1
config-pin P8_39 0

#SWITCH 2
config-pin P8_40 1
config-pin P8_40 0

#Switch 4
config-pin P8_42 1
config-pin P8_42 0

#Switch 3
config-pin P8_37 1
config-pin P8_37 0

#Switch 5
config-pin P8_45 1
config-pin P8_45 0

#Switch 6
config-pin P8_46 1
config-pin P8_46 0

#CAN1
config-pin p9.24 can
config-pin p9.26 can

#CAN0
config-pin P9_19 can
config-pin P9_20 can

exit 0
